# deploy/backup/suxios-db-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
25 3 * * * root /usr/local/sbin/suxios-db-backup >> /var/log/suxios-db-backup.log 2>&1
10 4 * * 0 root /usr/local/sbin/suxios-db-restore-verify >> /var/log/suxios-db-restore-verify.log 2>&1
