# deploy/systemd/suxios-cloud-hotel-daily@.service
[Unit]
Description=SUXIOS hotel-scoped daily automation (hotel %i)
After=network-online.target mariadb.service php8.3-fpm.service
Wants=network-online.target
ConditionPathExists=/var/www/suxios/current/think
StartLimitIntervalSec=2h
StartLimitBurst=10

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/suxios/current
Environment=TZ=Asia/Shanghai
EnvironmentFile=-/etc/suxios/cloud-automation.env
ExecStart=/usr/bin/php -d memory_limit=384M think cloud-automation:run --mode=daily --hotel-id=%i --limit=1 --max-attempts=10 --lock-wait-seconds=1500
Restart=on-failure
RestartSec=2min
Nice=10
CPUQuota=80%
MemoryMax=512M
TasksMax=32
TimeoutStartSec=30min
UMask=0027
