# deploy/systemd/suxios-cloud-hotel-daily@.timer
[Unit]
Description=Run one SUXIOS hotel's yesterday operating loop every morning
Conflicts=suxios-cloud-daily.timer

[Timer]
OnCalendar=*-*-* 08:10:00 Asia/Shanghai
RandomizedDelaySec=3min

[Install]
WantedBy=timers.target
